# files/xfce-binaries/install/usr/lib/systemd/user/xfconfd.service
[Unit]
Description=Xfce configuration service

[Service]
Type=dbus
BusName=org.xfce.Xfconf
ExecStart=/usr/lib64/xfce4/xfconf/xfconfd
